# systemd_files/eod_jobs.timer
[Unit]
Description=Run eod_jobs.service daily at 23:50 UTC

[Timer]
# Run every day at 23:50 UTC (explicit UTC so DST/local TZ don't affect it)
OnCalendar=*-*-* 23:50:00 UTC
Persistent=true
Unit=eod_jobs.service
# Optional jitter / precision:
# RandomizedDelaySec=2min
# AccuracySec=1s

[Install]
WantedBy=timers.target
