# systemd_files/live_pnl.service
[Unit]
Description=Write daily Betfair equity/PnL snapshot

[Service]
Type=oneshot
WorkingDirectory=/root
ExecStart=/root/betting/.venv/bin/python -u /root/betting/jobs/live_pnl_dump.py
