# systemd_files/live_pnl.service
[Unit]
Description=Write daily Betfair equity/PnL snapshot

[Service]
Type=oneshot
WorkingDirectory=/root
ExecStart=/root/betfair_profitbox/.venv/bin/python -u /root/betfair_profitbox/jobs/live_pnl_dump.py
