# systemd_files/md_recorder-restart.service
[Unit]
Description=Restart md_recorder.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart md_recorder.service
