# systemd_files/strat_liquidity_momentum.service
# /etc/systemd/system/strat_liquidity_momentum.service
[Unit]
Description=Always-on strat_liquidity_momentum (liquidity momentum strategy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
# EnvironmentFile=-/root/betfair_profitbox/.env
Environment="PYTHONUNBUFFERED=1" "PYTHONPATH=/root"
ExecStart=/root/betfair_profitbox/.venv/bin/python -u -m betfair_profitbox.jobs.strat_liquidity_momentum
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=journal
StandardError=inherit

[Install]
WantedBy=multi-user.target
