# systemd_files/strat_risk_backfave.service
# /etc/systemd/system/strat_risk_backfave.service
[Unit]
Description=Always-on strat_risk_backfave (risk backfave strategy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
# Drop the file dependency unless you actually use it
# EnvironmentFile=-/root/betfair_profitbox/.env
# Inline what you need:
Environment="PYTHONUNBUFFERED=1" "PYTHONPATH=/root"
# RUN AS MODULE so absolute imports work
ExecStart=/root/betfair_profitbox/.venv/bin/python -u -m betfair_profitbox.jobs.strat_risk_backfave
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=journal
StandardError=inherit

[Install]
WantedBy=multi-user.target

