# systemd_files/strat_trades_report.timer
[Unit]
Description=Run strat_trades_report.service daily at 22:30 UTC

[Timer]
# Run every day at 22:30 UTC (explicit UTC so DST/local TZ don't affect it)
OnCalendar=*-*-* 22:30:00 UTC
Persistent=true
Unit=strat_trades_report.service
# Optional jitter / precision:
# RandomizedDelaySec=2min
# AccuracySec=1s

[Install]
WantedBy=timers.target
